# scripts/saviaclaw-backup.service
[Unit]
Description=SaviaClaw Daily Memory Backup

[Service]
Type=oneshot
User=monica
ExecStart=/bin/bash /home/monica/claude/scripts/saviaclaw-backup.sh
