# scripts/saviaclaw-backup.timer
[Unit]
Description=SaviaClaw Daily Memory Backup Timer

[Timer]
OnCalendar=*-*-* 22:30:00
Persistent=true

[Install]
WantedBy=timers.target
